# templates/console/config/upstart/init.conf
# generated by puppet from serial/templates/console/config/upstart/init.conf
#
# <%= @name %> - agetty
#
# This service maintains a agetty on <%= name %>.

stop on runlevel [S06]
start on runlevel [1235]

respawn
exec agetty -h -L -w /dev/<%= @name %> <%= @baud %> vt102
